# scheduled-zfs-sync.service
[Unit]
Description=Scheduled ZFS Sync with RTC Wake and Tailscale Wait
Requires=network-online.target tailscaled.service
After=network-online.target tailscaled.service

[Service]
Type=oneshot
User=root
ExecStartPre=/bin/bash -c 'until tailscale status --json | jq -e ".Self.Online == true" >/dev/null; do sleep 5; done'
ExecStart=/usr/local/bin/scheduled-zfs-sync.sh
